# systemd/bms-background.service
[Unit]
Description=BMS Background Data Collection Service
Documentation=file:///home/xavier-bms/Documents/bms-native-reader/raspberry_pi_bms_reader/README.md
After=network-online.target bluetooth.service
Wants=network-online.target
Requires=bluetooth.service

[Service]
Type=simple
User=xavier-bms
Group=xavier-bms
WorkingDirectory=/home/xavier-bms/Documents/bms-native-reader/raspberry_pi_bms_reader
Environment=PATH=/home/xavier-bms/Documents/bms-native-reader/raspberry_pi_bms_reader/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/xavier-bms/Documents/bms-native-reader/raspberry_pi_bms_reader/venv/bin/python /home/xavier-bms/Documents/bms-native-reader/raspberry_pi_bms_reader/bms_background_service.py --data-file /tmp/bms_latest.json --interval 5.0 --log-level INFO
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bms-background

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp

# Resource limits
MemoryMax=256M
CPUQuota=50%

# Supplementary groups for Bluetooth access
SupplementaryGroups=bluetooth

[Install]
WantedBy=multi-user.target